# _OS_SETTINGS/etc/udev/rules.d/rune_usb-audio.rules
#KERNEL=="card0", DRIVER=="snd-usb-audio", RUN+="/var/www/command/refresh_ao"
KERNEL=="card0", SUBSYSTEM=="sound", RUN+="/var/www/command/refresh_ao"